# dsync-hourly.timer
[Unit]
Description=dsync hourly timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
