# src/dsync-hourly.service
[Unit]
Description=Run dsync hourly

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dsync
